# roles/services/templates/mkbackup.timer
[Unit]
Description=Make backup of data (timer)

[Timer]
Unit=mkbackup.service
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
